# contrib/daemontools/log/run
#!/bin/sh

exec \
 /usr/local/bin/setuidgid nobody
  /usr/local/bin/multilog t s10485760 n5 \
  /var/log/i2s3.log
